# rpm/fingertip-shared-cache.service
[Unit]
Description=Ensure fingertip shared cache is CoW-enabled

[Service]
Type=oneshot
ExecStart=/usr/libexec/fingertip/shared-cache-setup
RemainAfterExit=yes

[Install]
WantedBy=default.target
